# config/pkmn.ai-backup.timer
# systemctl enable /home/pkmn/ai/config/pkmn.ai-backup.timer

[Unit]
Description=Monthly archive, upload, and cleanup of old battle logs

[Timer]
OnCalendar=*-*-01 01:00:00

[Install]
WantedBy=timers.target